# safe-defaults/eos-safe-defaults.service.in
[Unit]
Description=Enable/update the Endless family-safe default settings
ConditionNeedsUpdate=/etc
Before=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@SBINDIR@/eos-safe-defaults enable

[Install]
WantedBy=multi-user.target
