# safe-defaults/eos-safe-defaults.service.in
[Unit]
Description=Enable/update the Endless family-safe default settings
DefaultDependencies=no
Requires=local-fs.target
After=local-fs.target
Before=NetworkManager.service systemd-update-done.service
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@SBINDIR@/eos-safe-defaults enable
ExecStop=@SBINDIR@/eos-safe-defaults disable

[Install]
WantedBy=multi-user.target
